# crons/health.cron
*/15 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /usr/bin/notify-send "Alongar pulso posterior"
*/10 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /usr/bin/notify-send "Alongar pulso"
*/20 * * * * XDG_RUNTIME_DIR=/run/user/$(id -u) /usr/bin/notify-send "Alongar pescoço e cervical"
15 * * * *   XDG_RUNTIME_DIR=/run/user/$(id -u) /usr/bin/notify-send "Alongar pernas"
