# crons/schedule_fix.cron
*/5 * * * * pkill -RTMIN+12 i3blocks
*/5 * * * * pkill -RTMIN+10 i3blocks
*/5 * * * * cat ~/petryfiles/crons/*.cron | crontab -
*/10 * * * * sudo systemctl restart vban.service
